# rootfs/etc/services.d/qdrant/run
#!/command/with-contenv bash
set -euo pipefail

exec /usr/local/bin/qdrant
