# rootfs/etc/s6-overlay/s6-rc.d/fastapi/run
#!/usr/bin/execlineb -P

# Export variables available to S6
with-contenv

# Change to api directory and boot Uvicorn
cd /app/api
uvicorn main:app --host 0.0.0.0 --port 8765 --workers 2
