# rootfs/etc/s6-overlay/s6-rc.d/nextjs/run
#!/usr/bin/execlineb -P

# Export variables available to S6
with-contenv

# Change to ui directory and start Next.js prod server
cd /app/ui
npm run start
