# rootfs/etc/s6-overlay/s6-rc.d/nextjs/run
#!/usr/bin/execlineb -P
cd /app/ui
npm run start
